# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
autostart=true
command=/usr/sbin/sshd -D

[program:mtm]
autostart=true
command=/usr/bin/mtm start
